# raspberry/servicios/safegrid-whatsapp.service
# /etc/systemd/system/safegrid-whatsapp.service
[Unit]
Description=SafeGrid WhatsApp Webhook (Gunicorn)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=teleadmin
WorkingDirectory=/home/teleadmin/safegrid-dashboard

EnvironmentFile=/home/teleadmin/safegrid.env
Environment=PYTHONUNBUFFERED=1

ExecStart=/home/teleadmin/venv-safegrid/bin/gunicorn -w 1 -b 0.0.0.0:3000 whatsapp_webhook:app

Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
