# raspberry/safegrid-dashboard/raspberry/servicios/safegrid-dashboard.service
# /etc/systemd/system/safegrid-dashboard.service
[Unit]
Description=SafeGrid Dashboard Web (Flask via Gunicorn)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=teleadmin
WorkingDirectory=/home/teleadmin/safegrid-dashboard
EnvironmentFile=/etc/safegrid.env

ExecStart=/home/teleadmin/venv-safegrid/bin/gunicorn -w 1 -b 0.0.0.0:8080 dashboard_server:app
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
