# raspberry/safegrid-dashboard/raspberry/servicios/safegrid-stats-writer.service
# /etc/systemd/system/safegrid-stats-writer.service
[Unit]
Description=SafeGrid Stats Writer (Pi+Net -> SQLite)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=teleadmin
WorkingDirectory=/home/teleadmin/safegrid-dashboard
EnvironmentFile=/etc/safegrid.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/teleadmin/venv-safegrid/bin/python /home/teleadmin/safegrid-dashboard/stats_writer.py
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
